# rc.local
#!/bin/sh -e
# Boblight START
/usr/local/bin/boblightd -f &
# Boblight END
exit 0